# recipes-perception/edgefirst-perception/files/edgefirst.target
[Unit]
Description=EdgeFirst Perception Services
After=multi-user.target
AllowIsolate=yes
Requires=multi-user.target
